# packages/audio/alsa-equal/system.d/alsa-equal.service
[Unit]
Description=Makes sure the default .asound.rc file is copied to the storage folder if it doesn't exist

[Service]
Type=oneshot
After=sound.target
ConditionPathExists=!$HOME/.asound.rc
ExecStart=/bin/cp /usr/share/alsa/.asound.rc $HOME/
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
