# packages/mediacenter/kodi/system.d/kodi-mysqlcheck.service
[Unit]
Description=Check for internet connection for possible Kodi mysql connection
ConditionPathExists=/storage/.cache/libreelec/network_wait
ConditionPathExists=/storage/.kodi/userdata/
Before=network-online.target
After=connman.service

[Service]
Type=oneshot
EnvironmentFile=/storage/.cache/libreelec/network_wait
ExecStart=/bin/sh -c '/usr/lib/systemd/system/kodi-mysqlcheck.sh ${WAIT_NETWORK_TIME}'
RemainAfterExit=yes

[Install]
WantedBy=network-online.target
